# launchagents/com.colima.lite-llm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.colima.lite-llm</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/renefernandoortegamunguia/.local/bin/start-colima-lite-llm</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/renefernandoortegamunguia/Library/Logs/colima-lite-llm.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/renefernandoortegamunguia/Library/Logs/colima-lite-llm.log</string>
  </dict>
</plist>
